# config/url.conf
[program:url]
command=/home/john/venv/bin/gunicorn volunteer_reports.wsgi:application -c /home/john/volunteer_reports/config/gunicorn.conf.py
directory=/home/john/volunteer_reports
user=john
autorestart=true
redirect_stderr=true
stdout_logfile=/home/john/volunteer_reports/logs/debug.log